# deploy_tools/gunicorn-upstart.template.conf
description "Gunicorn server for superlists-staging.ottg.eu"

start on net-device-up
stop on shutdown

respawn

setuid fox
chdir /home/fox/sites/superlists-staging.ottg.eu/source
			 
exec ../virtualenv/bin/gunicorn --bind unix:/tmp/superlists-staging.ottu.eu.socket superlists.wsgi:application